# _release/cacophony-processing.thermal@.service
[Unit]
Description=Cacophony Project thermal upload processing
After=network.target

[Service]
ExecStart=/usr/bin/cacophony-processing.pex -m thermal_processing
Restart=on-failure
User=fullnoise
Group=fullnoise

[Install]
WantedBy=multi-user.target
